# config/supervisord.conf
[supervisord]
nodaemon=true

[program:mosquitto]
command=/usr/sbin/mosquitto -c /etc/mosquitto/mosquitto.conf
autorestart=true

[program:go-server]
command=/app/server
autorestart=true
